# deploy/jira-summary.service
[Unit]
Description=Jira Summary Scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/jira-summary
ExecStart=/usr/bin/node scripts/scheduler.js
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
