# deploy/jira-summary.timer
[Unit]
Description=Run Jira Summary at 20:00 on weekdays (Asia/Ho_Chi_Minh)

[Timer]
# 20:00 GMT+7 = 13:00 UTC
OnCalendar=Mon..Fri 13:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
